# etc/systemd/user/bgworker.slice
[Unit]
Description=Background Worker Slice
Before=slices.target

[Slice]
MemoryLow=1G
MemoryHigh=2G
MemoryLimit=3G
CPUWeight=25
