# etc/systemd/user/bgworker.slice
[Unit]
Description=Background Worker Slice
Before=slices.target

[Slice]
MemoryLow=2G
CPUWeight=25
