# etc/systemd/user/browser.slice
[Unit]
Description=Browser Slice
Before=slices.target

[Slice]
MemoryLow=4G
CPUWeight=25
